# Dockerfiles/openvpn/container-root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv sh

exec openvpn /config/openvpn/openvpn.conf
