# Dockerfiles/rtorrent/container-root/etc/services.d/rtorrent/finish
#!/usr/bin/with-contenv sh

# Make sure lock file is removed after end
rm -f /config/rtorrent/rtorrent.lock
