# Dockerfiles/sonarr/services.d/sonarr/run
#!/usr/bin/with-contenv sh

mkdir -p /config/sonarr

rm -f /config/sonarr/nzbdrone.pid

exec mono /opt/NzbDrone/NzbDrone.exe -nobrowser -data=/config/sonarr
