# Dockerfiles/rtorrent/container-root/etc/services.d/php/run
#!/usr/bin/with-contenv sh

exec php-fpm7.0 -F -R
