# Dockerfiles/h5ai/supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command=/usr/sbin/php5-fpm --nodaemonize

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"