# Dockerfiles/sonarr/services.d/jackett/run
#!/usr/bin/with-contenv sh

exec mono /opt/Jackett/JackettConsole.exe
